# buildsteps/packages/018_prometheus-boinc-exporter/template/etc/service/boinc_exporter/finish
#!/bin/sh

source /usr/share/uebo_shell_stuff
deregister_prometheus_service boinc_exporter
