# buildsteps/packages/018_prometheus-boinc-exporter/template/etc/service/boinc_exporter/run
#!/bin/sh

exec 2>&1

if [ ! -e /var/lib/boinc-client/gui_rpc_auth.cfg ]; then
	exit 1
fi

source /usr/share/uebo_shell_stuff
register_prometheus_service boinc_exporter 9109

memlimit=$((40*1024*1024)) # 40MB?
create_memory_group prometheus $memlimit

exec cgexec -g memory:prometheus /usr/bin/boinc_exporter "$(cat /var/lib/boinc-client/gui_rpc_auth.cfg)"
