# buildsteps/rootfiles/etc/service.d/getty/run
#!/bin/sh

consolename=$(basename $(pwd))

exec /sbin/getty 0 $consolename
