# buildsteps/rootfiles/etc/service.d/getty/run
#!/bin/sh

consolename=$(basename $(pwd) | cut -d"-" -f2-)

exec /sbin/getty 0 $consolename
