# buildsteps/rootfiles/etc/service.d/dhclient/run
#!/bin/sh

ifname=$(basename $(pwd))

exec /sbin/udhcpc -f -i "${ifname}"
