# buildsteps/rootfiles/etc/service.d/dhclient/run
#!/bin/sh

ifname=$(basename $(pwd) | cut -d"-" -f2-)

exec /sbin/udhcpc -f -i "${ifname}"
