# buildsteps/rootfiles/etc/service/node_exporter/finish
#!/bin/sh

source /usr/share/uebo_shell_stuff

deregister_prometheus_service node_exporter
