# buildsteps/rootfiles/etc/service/dropbear/run
#!/bin/sh

exec /usr/sbin/dropbear -R -F -g
