# buildsteps/rootfiles/etc/service/ntpd/run
#!/bin/sh

exec 2>&1

peers=""

if [ ! -f "/etc/ntp.conf" ]; then
	peers="-p pool.ntp.org"
else
	source /etc/ntp.conf
fi

exec /usr/sbin/ntpd -n ${peers}
